# apps/api/deploy/supervisor/belsuite-queue-workers.conf
; Supervisor programs for Laravel queue workers (Octane-compatible: separate
; processes from HTTP workers; restart after max-time to cap memory).
;
; Install: copy to /etc/supervisor/conf.d/ and `supervisorctl reread && supervisorctl update`
; Adjust paths and user for your server.

[program:belsuite-queue-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/belsuite/apps/api/artisan queue:work failover_redis --sleep=1 --rest=1 --queue=default --tries=3 --max-time=3600 --memory=256
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/belsuite/apps/api/storage/logs/worker-default.log
stopwaitsecs=3600

[program:belsuite-queue-ai]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/belsuite/apps/api/artisan queue:work failover_redis --sleep=1 --rest=1 --queue=ai --tries=4 --timeout=700 --max-time=1800 --memory=512
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/belsuite/apps/api/storage/logs/worker-ai.log
stopwaitsecs=7200

[program:belsuite-queue-video]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/belsuite/apps/api/artisan queue:work failover_redis --sleep=2 --rest=2 --queue=video --tries=3 --timeout=920 --max-time=7200 --memory=1024
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/belsuite/apps/api/storage/logs/worker-video.log
stopwaitsecs=7200

[program:belsuite-queue-campaigns]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/belsuite/apps/api/artisan queue:work failover_redis --sleep=1 --rest=1 --queue=campaigns --tries=5 --timeout=500 --max-time=3600 --memory=256
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/belsuite/apps/api/storage/logs/worker-campaigns.log
stopwaitsecs=3600
